# tradingalert.service
[Unit]
Description=Trading Alert Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tradingalert
# bind to localhost:5000, Nginx proxies to it
ExecStart=/opt/tradingalert/venv/bin/gunicorn --workers 1 --bind 127.0.0.1:5000 "webapp:create_app()"
Restart=always
RestartSec=60
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
